# network-restart/network-restart-actual-restart.service
[Unit]
Description=Reboot Device

[Service]
Type=oneshot
ExecStart=/usr/bin/reboot 0