# network-restart/network-restart.timer
[Unit]
Description=Check for network connectivity every 10mins

[Timer]
OnUnitActiveSec=10m
OnBootSec=10m

[Install]
WantedBy=timers.target
