# restart-wifi/restart-wifi-actual-run.service
[Unit]
Description=Restart wifi adapter

[Service]
Type=oneshot
ExecStart=/usr/bin/restart-wifi
