# restart-wifi/restart-wifi-check.timer
[Unit]
Description=Check for network connectivity every 10mins

[Timer]
OnUnitActiveSec=1m
OnBootSec=1m

[Install]
WantedBy=timers.target
