# laptop-low-battery/laptop-low-battery.service
[Unit]
Description=Check for Low Battery

[Service]
Type=oneshot
ExecStart=/usr/bin/laptop-low-battery
